# roles/gateway/templates/gateway-direct-mode.service.j2
[Unit]
Description=Gateway Direct Mode Service
Documentation=man:systemd.service(5)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gateway-apply-rules direct
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
