# roles/backup/templates/backup.timer.j2
[Unit]
Description=Backup {{ backup_name }} timer

[Timer]
OnCalendar={{ backup_systemd_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
