# roles/cloudflared/templates/cloudflared.service.j2
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate run --token {{ cloudflared_tunnel_token }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
